# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port = :9001
username = admin
password = changeme

[program:sshd]
command=/usr/sbin/sshd -D
